# other/led_clock.service
[Unit]
Description=Led clock
After=lircd.service

[Service]
ExecStart=/opt/led_clock/led_clock.py

[Install]
WantedBy=multi-user.target
